# curator4-1695918829.conf
[program:curator4-1695918829]
command=python3 /home/curator4/curator4-1695918829/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-1695918829/test.err.log
stdout_logfile=/home/curator4/curator4-1695918829/test.out.log